# eve/workers/orbit/services/push-server/run
#!/usr/bin/with-contenv sh

cd /go/src/${PACKAGE}

exec ./pensieve-push-server
