# eve/workers/orbit/services/configurator/run
#!/usr/bin/with-contenv sh

config_mongo.sh

touch /tmp/hold_open
tail -f /tmp/hold_open
