# eve/workers/orbit/services/mongod/run
#!/usr/bin/with-contenv sh

mkdir -p /tmp/mongod

exec mongod --dbpath /tmp/mongod
