# tests/sfko/s6/services.d/sfko/run
#!/usr/bin/with-contenv sh

cd /usr/bin/sfko
su-exec sfko python run.py
